# linux/systemd/git-autopull.timer
[Unit]
Description=Poll Codeberg for homelab repo changes every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
