# linux/systemd/pmtiles-rebuild.service
[Unit]
Description=Mirror latest OpenFreeMap planet PMTiles (download MBTiles, verify SHA256, convert with go-pmtiles)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=webdev
Group=webdev
ExecStart=/usr/local/sbin/pmtiles-rebuild.sh
TimeoutStartSec=12h
MemoryMax=8G
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/webdev/tiles
PrivateTmp=true
NoNewPrivileges=true
