# linux/systemd/pmtiles-rebuild.timer
[Unit]
Description=Monthly OpenFreeMap planet PMTiles mirror refresh
Requires=pmtiles-rebuild.service

[Timer]
# OpenFreeMap publishes weekly; once a month is plenty for our consumer (lovelang).
OnCalendar=*-*-08 04:00:00
RandomizedDelaySec=2h
Persistent=true

[Install]
WantedBy=timers.target
