# linux/systemd/pmtiles-rebuild.service
[Unit]
Description=Rebuild Ireland PMTiles bundle from latest Geofabrik OSM PBF via tilemaker (OpenMapTiles schema)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=webdev
Group=webdev
ExecStart=/usr/local/sbin/pmtiles-rebuild.sh
TimeoutStartSec=60min
MemoryMax=8G
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/webdev/tiles
PrivateTmp=true
NoNewPrivileges=true
