# linux/systemd/trivy-scan.timer
[Unit]
Description=Daily Homelab CVE scan (Trivy)

[Timer]
OnCalendar=*-*-* 03:17:00
RandomizedDelaySec=30m
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
