# launchagents/com.noc.teamspeak.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noc.teamspeak</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/noc/teamspeak3-server_mac/ts3server</string>
        <string>inifile=ts3server.ini</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/noc/teamspeak3-server_mac</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/noc/noc-homelab/logs/teamspeak.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/noc/noc-homelab/logs/teamspeak-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TS3SERVER_LICENSE</key>
        <string>accept</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
