# launchagents/com.noc.copyparty.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noc.copyparty</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>: > /Users/noc/Library/Logs/noc-homelab/copyparty.log; : > /Users/noc/Library/Logs/noc-homelab/copyparty.error.log; exec /Users/noc/.local/bin/copyparty -p 8081 -i 0.0.0.0 -v /Users:/users:rw</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/noc/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/copyparty.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/copyparty.error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/noc</string>
</dict>
</plist>
