# launchagents/com.noc.mdsf-crew-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noc.mdsf-crew-api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>: > /Users/noc/Library/Logs/noc-homelab/mdsf-crew-api.log; : > /Users/noc/Library/Logs/noc-homelab/mdsf-crew-api.error.log; exec /opt/homebrew/bin/tsx watch src/index.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/noc/dev/mdsf-crew/apps/api</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/mdsf-crew-api.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/mdsf-crew-api.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/noc/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/noc</string>
        <key>NODE_ENV</key>
        <string>development</string>
    </dict>
</dict>
</plist>
