# launchagents/com.noc.mdsf-org.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noc.mdsf-org</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>: &gt; /Users/noc/Library/Logs/noc-homelab/mdsf-org.log; : &gt; /Users/noc/Library/Logs/noc-homelab/mdsf-org.error.log; exec /opt/homebrew/bin/python3 -m http.server 8190 --directory /Users/noc/dev/mdsf-core &gt;&gt; /Users/noc/Library/Logs/noc-homelab/mdsf-org.log 2&gt;&gt; /Users/noc/Library/Logs/noc-homelab/mdsf-org.error.log</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/noc/dev/mdsf-core</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/mdsf-org.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/mdsf-org.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/noc/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/noc</string>
    </dict>
</dict>
</plist>
