# launchagents/com.noc.homelab-status-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noc.homelab-status-push</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/noc/noc-homelab/scripts/push-homelab-status.py</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/homelab-status-push.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/homelab-status-push.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/noc</string>
        <key>HOMELAB_STATUS_WORKER_URL</key>
        <string>https://homelab-status.nocturnef4.workers.dev</string>
        <key>HOMELAB_STATUS_PUSH_SECRET</key>
        <string>KkIY9XKSiNfQJ4pg8YKRnjk_y1EEkfFgzlGwH6ME0Ss</string>
    </dict>
</dict>
</plist>
