# launchagents/com.maloja.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.maloja.service</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/noc/.local/bin/maloja</string>
        <string>run</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/noc</string>
    <key>StandardOutPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/maloja.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/maloja.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/noc/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/noc</string>
        <key>MALOJA_SKIP_SETUP</key>
        <string>yes</string>
        <key>MALOJA_FORCE_PASSWORD</key>
        <string>changeme123</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
