# launchagents/com.noc.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noc.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>: > /Users/noc/Library/Logs/noc-homelab/dashboard.log; : > /Users/noc/Library/Logs/noc-homelab/dashboard.error.log; exec /opt/homebrew/bin/python3 /Users/noc/noc-homelab/dashboard/app.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/noc/noc-homelab/dashboard</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/dashboard.error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/noc</string>
    </dict>
</dict>
</plist>
