# launchagents/com.noc.forgejo-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noc.forgejo-runner</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/noc/forgejo/runners/forgejo-runner</string>
        <string>daemon</string>
        <string>--config</string>
        <string>/Users/noc/forgejo/runners/config.yaml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/noc/forgejo/runners</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/forgejo-runner.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/forgejo-runner.error.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
