# launchagents/com.noc.dns-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.noc.dns-audit</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/noc/noc-homelab/scripts/dns-audit.sh</string>
        <string>--quiet</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/dns-audit.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/noc/Library/Logs/noc-homelab/dns-audit.error.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
